# files/hs5t.init
#!/bin/sh /etc/rc.common

. /lib/functions.sh

START=90
USE_PROCD=1

start_hs5t() {
	procd_open_instance
	procd_set_param respawn
	procd_set_param stderr 1
	procd_set_param command /usr/bin/hs5t
	procd_append_param command /etc/hs5t.yml
	procd_close_instance
}

start_service() {
	start_hs5t
}
